# installer/centos/sysmanage.service
[Unit]
Description=SysManage Server
Documentation=https://sysmanage.org
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=sysmanage
Group=sysmanage
WorkingDirectory=/opt/sysmanage
Environment="SYSMANAGE_CONFIG=/etc/sysmanage.yaml"
Environment="SYSMANAGE_DB_DIR=/var/lib/sysmanage"
Environment="SYSMANAGE_LOG_DIR=/var/log/sysmanage"
ExecStart=/opt/sysmanage/.venv/bin/python -m uvicorn backend.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sysmanage

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sysmanage /var/log/sysmanage /etc/sysmanage

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
