# installer/ubuntu/debian/postinst
#!/bin/bash
set -e

# Automatically added by dh_installinit/13.11.4
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# Create sysmanage user if it doesn't exist
	if ! getent passwd sysmanage >/dev/null; then
		echo "Creating sysmanage system user..."
		adduser --system --group --home /nonexistent --no-create-home \
			--disabled-password --disabled-login \
			--gecos "SysManage Server" sysmanage
	fi

	# Create required directories
	mkdir -p /var/lib/sysmanage
	mkdir -p /var/log/sysmanage
	mkdir -p /etc/sysmanage

	# Set ownership of application directories
	chown -R sysmanage:sysmanage /opt/sysmanage
	chown -R sysmanage:sysmanage /var/lib/sysmanage
	chown -R sysmanage:sysmanage /var/log/sysmanage
	chown -R sysmanage:sysmanage /etc/sysmanage

	# Set proper permissions
	chmod 755 /opt/sysmanage
	chmod 755 /var/lib/sysmanage
	chmod 755 /var/log/sysmanage
	chmod 750 /etc/sysmanage

	# Fix virtualenv to work with system Python
	echo "Setting up Python virtual environment..."
	cd /opt/sysmanage
	rm -rf .venv
	python3 -m venv .venv
	.venv/bin/pip install --quiet --upgrade pip
	# Install runtime-only deps.  The full requirements.txt carries dev
	# tooling (pylint/astroid/pytest/playwright) that the server does not
	# need to run AND that the air-gap wheel set does not contain — using
	# it offline fails with "No matching distribution found for astroid".
	# Fall back to requirements.txt only if the prod file isn't shipped.
	if [ -f requirements-prod.txt ]; then
		.venv/bin/pip install --quiet -r requirements-prod.txt
	else
		.venv/bin/pip install --quiet -r requirements.txt
	fi
	cd -
	echo "✓ Python dependencies installed"

	# Create config file if it doesn't exist
	if [ ! -f /etc/sysmanage.yaml ]; then
		echo "Creating default configuration file..."
		cp /etc/sysmanage/sysmanage.yaml.example /etc/sysmanage.yaml
		chown sysmanage:sysmanage /etc/sysmanage.yaml
		chmod 640 /etc/sysmanage.yaml
	fi

	# ---------------------------------------------------------------
	# OpenBAO (secrets broker) — install + start + initialize/unseal.
	#
	# OpenBAO is central to SysManage.  We install the native package
	# (it provides the `bao` binary, the `openbao` service user, and the
	# stock openbao.service unit), drop our own minimal config, and run a
	# one-shot init/unseal so a fresh box comes up with a working local
	# secrets broker.  See docs/planning/openbao-deployment-and-airgap.md.
	# ---------------------------------------------------------------
	if ! command -v bao >/dev/null 2>&1; then
		echo "Installing OpenBAO..."
		# 1) Air-gap path: a bundled .deb staged next to the SysManage
		#    package (the air-gap mega-ISO places it here).
		BUNDLED_BAO="$(ls /opt/sysmanage/installer/openbao/*.deb 2>/dev/null | head -n1 || true)"
		if [ -n "$BUNDLED_BAO" ]; then
			apt-get install -y "$BUNDLED_BAO" >/dev/null 2>&1 \
				|| dpkg -i "$BUNDLED_BAO" >/dev/null 2>&1 || true
		# 2) Online path: configure the official OpenBAO apt repo, install.
		elif command -v curl >/dev/null 2>&1; then
			install -m0755 -d /etc/apt/keyrings 2>/dev/null || true
			if curl -fsSL https://pkg.openbao.org/gpg.key 2>/dev/null \
				| gpg --dearmor -o /etc/apt/keyrings/openbao.gpg 2>/dev/null; then
				echo "deb [signed-by=/etc/apt/keyrings/openbao.gpg] https://pkg.openbao.org/deb/ubuntu stable main" \
					> /etc/apt/sources.list.d/openbao.list
				apt-get update >/dev/null 2>&1 || true
				apt-get install -y openbao >/dev/null 2>&1 || true
			fi
		fi
		if ! command -v bao >/dev/null 2>&1; then
			echo "⚠ OpenBAO ('bao') could not be installed automatically."
			echo "  Install it manually (https://openbao.org) or set vault.enabled=false."
		fi
	fi

	# Provision OpenBAO directories + config when the binary is present.
	if command -v bao >/dev/null 2>&1; then
		# The openbao package creates the 'openbao' user; fall back to
		# creating it if a tarball-based install didn't.
		if ! getent passwd openbao >/dev/null; then
			adduser --system --group --home /var/lib/openbao --no-create-home \
				--disabled-password --disabled-login \
				--gecos "OpenBAO" openbao >/dev/null 2>&1 || true
		fi
		mkdir -p /var/lib/openbao/data /etc/openbao
		chown -R openbao:openbao /var/lib/openbao
		# Our config ships to /etc/openbao/openbao.hcl via debian/install.
		chown root:openbao /etc/openbao/openbao.hcl 2>/dev/null || true
		chmod 640 /etc/openbao/openbao.hcl 2>/dev/null || true
	fi

	# Reload systemd daemon
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null 2>&1 || true
	fi

	# Enable + start OpenBAO and run the init/unseal one-shot so the local
	# secrets broker is ready before the operator starts SysManage.  Safe to
	# run now (OpenBAO needs no external config); SysManage itself is still
	# left stopped for the operator to configure.
	if command -v bao >/dev/null 2>&1 && [ -d /run/systemd/system ]; then
		systemctl enable openbao.service >/dev/null 2>&1 || true
		systemctl enable sysmanage-openbao-init.service >/dev/null 2>&1 || true
		systemctl restart openbao.service >/dev/null 2>&1 || true
		systemctl start sysmanage-openbao-init.service >/dev/null 2>&1 \
			|| echo "⚠ OpenBAO init/unseal did not complete; check 'systemctl status sysmanage-openbao-init'."
	fi

	# Enable but DON'T start the service (user needs to configure first)
	if [ -x "/usr/bin/deb-systemd-invoke" ]; then
		deb-systemd-invoke enable sysmanage.service >/dev/null 2>&1 || true
	elif [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper enable sysmanage.service >/dev/null 2>&1 || true
	else
		systemctl enable sysmanage.service >/dev/null 2>&1 || true
	fi

	# Configure nginx if installed
	if [ -x "/usr/sbin/nginx" ]; then
		echo "Configuring nginx..."
		# Enable the sysmanage site
		if [ -f /etc/nginx/sites-available/sysmanage-nginx.conf ]; then
			ln -sf /etc/nginx/sites-available/sysmanage-nginx.conf /etc/nginx/sites-enabled/sysmanage 2>/dev/null || true
			# Test nginx config
			nginx -t >/dev/null 2>&1 && systemctl reload nginx >/dev/null 2>&1 || echo "⚠ nginx configuration may need manual review"
		fi
	fi

	# On a configure, $2 holds the previously-installed version: empty on a
	# fresh install, set on an in-place upgrade.  The upgrade path leaves the
	# existing config + database in place, so print the short migrate+restart
	# checklist instead of the full first-time configuration wall.
	if [ -n "$2" ]; then
	echo ""
	echo "=========================================="
	echo "SysManage Server upgrade complete!"
	echo "=========================================="
	echo ""
	echo "Your existing /etc/sysmanage.yaml was preserved and the service was"
	echo "stopped for the upgrade.  It is NOT restarted automatically — two"
	echo "manual steps remain:"
	echo ""
	echo "1. Apply database migrations (all three partition chains)"
	echo "   ------------------------"
	echo "     cd /opt/sysmanage"
	echo "     sudo -u sysmanage .venv/bin/python -m alembic --name registry upgrade head"
	echo "     sudo -u sysmanage .venv/bin/python -m alembic --name shared upgrade head"
	echo "     sudo -u sysmanage .venv/bin/python -m alembic upgrade head"
	echo ""
	echo "2. Restart the service"
	echo "   -------------------"
	echo "     sudo systemctl restart sysmanage"
	echo "     sudo systemctl status sysmanage"
	echo ""
	echo "Log files: /var/log/sysmanage/"
	echo "Web interface: https://your-server:8443/"
	echo ""
	else
	echo ""
	echo "=========================================="
	echo "SysManage Server installation complete!"
	echo "=========================================="
	echo ""
	echo "⚠ IMPORTANT: Configuration Required"
	echo ""
	echo "Before starting SysManage, you MUST configure:"
	echo ""
	echo "1. PostgreSQL Database Connection"
	echo "   --------------------------------"
	echo "   SysManage requires a PostgreSQL database (version 12+)."
	echo "   The database can be on this server or a remote host."
	echo ""
	echo "   To set up a local database:"
	echo "     sudo apt-get install postgresql"
	echo "     sudo -u postgres createuser sysmanage"
	echo "     sudo -u postgres createdb sysmanage -O sysmanage"
	echo "     sudo -u postgres psql -c \"ALTER USER sysmanage WITH PASSWORD 'your-password';\""
	echo ""
	echo "2. Configuration File"
	echo "   ------------------"
	echo "   Edit: /etc/sysmanage.yaml"
	echo ""
	echo "   Use the online configuration builder at:"
	echo "   https://sysmanage.org/config-builder.html"
	echo ""
	echo "   This tool will help you generate a complete configuration"
	echo "   with all required settings including:"
	echo "     • Database connection (PostgreSQL)"
	echo "     • Server hostname and ports"
	echo "     • SSL/TLS certificates"
	echo "     • Authentication settings"
	echo "     • Logging configuration"
	echo ""
	echo "3. Database Migration (all three partition chains)"
	echo "   ------------------"
	echo "   After configuring the database connection, run:"
	echo "     cd /opt/sysmanage"
	echo "     sudo -u sysmanage .venv/bin/python -m alembic --name registry upgrade head"
	echo "     sudo -u sysmanage .venv/bin/python -m alembic --name shared upgrade head"
	echo "     sudo -u sysmanage .venv/bin/python -m alembic upgrade head"
	echo ""
	echo "4. Start the Service"
	echo "   -----------------"
	echo "   After configuration and migration:"
	echo "     sudo systemctl start sysmanage"
	echo "     sudo systemctl status sysmanage"
	echo ""
	echo "Log files: /var/log/sysmanage/"
	echo "Web interface: https://your-server:8443/ (after configuration)"
	echo ""
	fi
fi

#DEBHELPER#

exit 0
